# supervisord.templates.d/04-rest-proxy.conf
[program:rest-proxy]
user=nobody
;environment=JMX_PORT=9583,KAFKAREST_JMX_OPTS="-Dcom.sun.management.jmxremote -Dcom.sun.management.jmxremote.authenticate=false -Dcom.sun.management.jmxremote.ssl=false -Dcom.sun.management.jmxremote.local.only=false -Djava.rmi.server.hostname=localhost -Dcom.sun.management.jmxremote.rmi.port=9583"
command=bash -c 'sleep 15 && exec /opt/landoop/kafka/bin/kafka-rest-start /opt/landoop/kafka/etc/kafka-rest/kafka-rest.properties'
redirect_stderr=true
stdout_logfile=/var/log/rest-proxy.log
startretries=5
