# supervisord.templates.d/09-nullsink-connector.conf
[program:nullsink]
user=nobody
command=bash -c 'sleep 180 && /usr/local/bin/nullsink.sh'
redirect_stderr=true
stdout_logfile=/var/log/nullsink.log
autorestart=false
startretries=1
