# filesystem/usr/local/share/landoop/etc/supervisord.templates.d/00-delayed-message.conf
[program:delayed-message]
command=bash -c 'for i in 4; do sleep $i; echo -e "$(envsubst < /usr/share/landoop/message)" > /proc/1/fd/1; done'
redirect_stderr=true
stdout_logfile=/dev/null
