# filesystem/usr/local/share/landoop/etc/supervisord.templates.d/06-caddy.conf
[program:caddy]
command=/opt/caddy/caddy -conf /var/run/caddy/Caddyfile
redirect_stderr=true
stdout_logfile=/var/log/caddy.log
