# filesystem/usr/local/share/landoop/etc/supervisord.templates.d/07-smoke-tests.conf
[program:smoke-tests]
directory=/tmp
command=/usr/local/bin/smoke-tests.sh
redirect_stderr=true
stdout_logfile=/var/log/smoke-tests.log
autorestart=false
startretries=1
