# filesystem/usr/local/share/landoop/etc/supervisord.templates.d/08-logs-to-kafka.conf
[program:logs-to-kafka]
command=bash -c 'sleep 180 && /usr/local/bin/logs-to-kafka.sh'
redirect_stderr=true
stdout_logfile=/var/log/logs-to-kafka.log
autorestart=false
startretries=1
