# filesystem/usr/local/share/landoop/etc/supervisord.templates.d/09-create-topics.conf
[program:create-topics]
user=nobody
command=bash -c 'eval $WAIT_SCRIPT_REGISTRY; exec /usr/local/bin/create-topics.sh'
redirect_stderr=true
stdout_logfile=/var/log/create-topics.log
startretries=5
