# filesystem/usr/local/share/landoop/etc/supervisord.templates.d/10-lenses.conf
[program:lenses]
user=nobody
environment=LENSES_HEAP_OPTS="-Xmx1288M -Xms250M"
directory=/opt/lenses
command=bash -c 'for ((i=0;i<15;i++)); do sleep 2; curl "http://localhost:8081" | grep "{}" && break; done; exec /opt/lenses/bin/lenses'
redirect_stderr=true
stdout_logfile=/var/log/lenses.log
