# supervisord.d/05-connect-distributed.conf
[program:connect-distributed]
user=nobody
; VCON is an unused variable, that we use to differentiate between the broker environment in a sed statement
environment=VCON=1,KAFKA_HEAP_OPTS=-Xmx{{CONNECT_HEAP}},KAFKA_JMX_OPTS="-Dcom.sun.management.jmxremote -Dcom.sun.management.jmxremote.authenticate=false -Dcom.sun.management.jmxremote.ssl=false -Dcom.sun.management.jmxremote.local.only=false -Djava.rmi.server.hostname=localhost -Dcom.sun.management.jmxremote.rmi.port=9584"
command=/opt/confluent/bin/connect-distributed /opt/confluent/etc/schema-registry/connect-avro-distributed.properties
redirect_stderr=true
stdout_logfile=/var/log/connect-distributed.log
startretries=5
