# supervisord.d/06-caddy.conf
[program:caddy]
user=nobody
command=/opt/caddy/caddy -conf /usr/share/landoop/Caddyfile
redirect_stderr=true
stdout_logfile=/var/log/caddy.log
