# extras/99-supervisord-running-sample-data.conf
[program:running-sample-data-ais]
user=nobody
command=bash -c 'sleep 45 && /usr/share/landoop/sample-data/running-ais.sh'
redirect_stderr=true
stdout_logfile=/var/log/running-ais.log
autorestart=false
startretries=1

[program:running-sample-data-taxis]
user=nobody
command=bash -c 'sleep 45 && /usr/share/landoop/sample-data/running-taxis.sh'
redirect_stderr=true
stdout_logfile=/var/log/running-taxis.log
autorestart=false
startretries=1

[program:running-sample-data-reddit]
user=nobody
command=bash -c 'sleep 45 && /usr/share/landoop/sample-data/running-reddit.sh'
redirect_stderr=true
stdout_logfile=/var/log/running-reddit.log
autorestart=false
startretries=1

[program:running-sample-data-backblaze-smart]
user=nobody
command=bash -c 'sleep 45 && /usr/share/landoop/sample-data/running-smart.sh'
redirect_stderr=true
stdout_logfile=/var/log/running-smart.log
autorestart=false
startretries=1

[program:running-sample-data-telecom-italia]
user=nobody
command=bash -c 'sleep 45 && /usr/share/landoop/sample-data/running-telecom-italia.sh'
redirect_stderr=true
stdout_logfile=/var/log/running-telecom-italia.log
autorestart=false
startretries=1

[program:running-cc-data]
user=nobody
environment=JAVA_HOME="/opt/lenses/jre8u131",JAVA_OPTS="-Xmx50m"
command=bash -c 'sleep 45 && /opt/generator/bin/generator /opt/generator/lenses.conf cc_data 1'
redirect_stderr=true
stdout_logfile=/var/log/running-cc-data.log
autorestart=false
startretries=1

[program:running-cc-payments]
user=nobody
environment=JAVA_HOME="/opt/lenses/jre8u131",JAVA_OPTS="-Xmx50m"
command=bash -c 'sleep 50 && /opt/generator/bin/generator /opt/generator/lenses.conf cc_payments 2 200'
redirect_stderr=true
stdout_logfile=/var/log/running-cc-payments.log
autorestart=false
startretries=1
