# extras/99-supervisord-sample-data.conf
[program:sample-data]
;user=nobody
command=bash -c 'sleep 45 && /usr/share/landoop/sample-data/put.sh'
redirect_stderr=true
stdout_logfile=/var/log/sampledata.log
autorestart=false
startretries=1
